# templates/default/cycle_server.init.erb
#!/bin/bash
#
# Generated by Chef for <%= node[:fqdn] %>
#
# Any local changes will be deleted.
#
# Startup script for Grill
# chkconfig: 35 86 15
# description: CycleServer (Tomcat and Derby)
# processname: cycle_server

# Grill's server.sh script provides init-script functionality.
# We just make sure you run it as the correct user.

# Source function library.
. /etc/rc.d/init.d/functions

DIR=`dirname $0`
CS_HOME=`/opt/cycle_server/system/scripts/get_home.sh`

export FULL_PATH="$0"
su - <%= @grill_user %>  -c "$CS_HOME/system/scripts/server.sh $@"


